# python-proxy.service
[Unit]
Description=Mastermind Python Proxy Suite v2.0
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/mastermind/protocols
ExecStart=/usr/bin/python3 /opt/mastermind/protocols/python_proxy.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment variables for proxy configuration
Environment=SOCKS_PORT=1080
Environment=WEBSOCKET_PORT=8080
Environment=HTTP_PROXY_PORT=8888
Environment=RESPONSE_PORTS=9000,9001,9002,9003
Environment=LOG_LEVEL=INFO
Environment=ENABLE_WEBSOCKET=true
Environment=ENABLE_HTTP_PROXY=true

# SSH configuration for WebSocket-to-SSH proxy
Environment=SSH_HOST=localhost
Environment=SSH_PORT=22
Environment=SSH_USER=root
Environment=SSH_PASS=
Environment=SSH_KEY_PATH=/root/.ssh/id_rsa

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/mastermind /opt/mastermind
PrivateTmp=true

[Install]
WantedBy=multi-user.target